# com.kangtearoom.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kangtearoom.update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/henry/telegram-stats/update.sh</string>
  </array>

  <!-- 매시 정각(하루 24회) + KST 08:55(직전 날 24시간 마감 수집). Minute만 주면 매시간 실행됨 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>55</integer></dict>
  </array>

  <!-- 노트북은 예약 시각에 꺼져/잠들어 있을 수 있다. launchd는 꺼진 Mac을 깨우지 못하므로
       RunAtLoad=true 로 로그인/부팅 직후 1회 실행해 놓친 회차를 보충한다.
       (잠자기 중 놓친 회차는 깨어날 때 launchd가 자동으로 1회 따라잡음) -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/henry/telegram-stats/data/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/henry/telegram-stats/data/launchd.log</string>
</dict>
</plist>
